# dot_config/systemd/user/mako.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/mako -c ${XDG_CONFIG_HOME}/niri/mako/config
